# UploadFiles/init-memtest.service
[Unit]
Description = Mem Test

[Service]
ExecStart = /boot/grub/initrdMemTest &

[Install]
WantedBy= multi-user.target
